# systemd/ima-research-bot.service
[Unit]
Description=IMA Research Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ima-search-bot
EnvironmentFile=/opt/ima-search-bot/.env
ExecStart=/opt/ima-search-bot/.venv/bin/python -m ima_research_bot serve
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
