# systemd/ima-research-radar.service
[Unit]
Description=IMA Research Radar
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/opt/ima-research-bot
EnvironmentFile=/opt/ima-research-bot/.env
ExecStart=/opt/ima-research-bot/.venv/bin/python -m ima_research_bot radar
